# root/etc/s6-overlay/s6-rc.d/svc-varroa/run
#!/command/execlineb -P

# - https://danyspin97.org/blog/getting-started-with-execline-scripting/
# - https://github.com/just-containers/s6-overlay

with-contenv
importas lsio_non_root_user LSIO_NON_ROOT_USER
importas wdir HOME

# If non-root user is set, run 's6-setuidgid'.
ifelse { test -z ${lsio_non_root_user} }
{
  # s6-notifyoncheck needs to run in the same folder as the service!
  # Check 'usage' section from here: https://skarnet.org/software/s6/s6-notifyoncheck.html
  s6-notifyoncheck -d -n 300 -w 1000
    foreground { echo "### Non-root mode ###" }
    s6-setuidgid abc
    cd ${wdir}
    # fdmove -c 3 1
    # fdmove -c 4 2
    # fdmove -c 2 1
    # redirfd -w 1 /config/svc-varroa.log
    foreground { echo "### Calling varroa ###" }
    ./varroa start
}
s6-notifyoncheck -d -n 300 -w 1000
  foreground { echo "### Root mode ###" }
  cd ${wdir}
  foreground { echo "### Calling varroa ###" }
  ./varroa start
    